# storage_pi_four/worker-boot.service
[Unit]
Description=Polar Eyes Worker Triage Service
# This ensures the script runs ONLY after the RAID is mounted.
# Systemd reads /etc/fstab and creates a .mount unit automatically.
Requires=mnt-raid.mount
After=mnt-raid.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/worker_triage.sh

[Install]
WantedBy=multi-user.target